# .config/systemd/user/ergodox.service
[Service]
ExecStart=/home/yuya373/dotfiles/restart_xkeysnail_service
